# flight-monitor-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <Triggers>
    <CalendarTrigger><StartBoundary>2026-05-21T08:00:00+08:00</StartBoundary><ScheduleByDay><DaysInterval>1</DaysInterval></ScheduleByDay></CalendarTrigger>
    <CalendarTrigger><StartBoundary>2026-05-21T12:00:00+08:00</StartBoundary><ScheduleByDay><DaysInterval>1</DaysInterval></ScheduleByDay></CalendarTrigger>
    <CalendarTrigger><StartBoundary>2026-05-21T20:00:00+08:00</StartBoundary><ScheduleByDay><DaysInterval>1</DaysInterval></ScheduleByDay></CalendarTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>python</Command>
      <Arguments>E:\surveillance-flight\surveillance.py --email</Arguments>
      <WorkingDirectory>E:\surveillance-flight</WorkingDirectory>
    </Exec>
  </Actions>
  <Settings><DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries></Settings>
</Task>
